# systemd/ili9488-daemon.service
[Unit]
Description=ILI9488 Display DMA SPI shared-memory daemon
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
User=pi
Group=pi
SupplementaryGroups=video spi gpio
AmbientCapabilities=CAP_SYS_RAWIO
CapabilityBoundingSet=CAP_SYS_RAWIO
EnvironmentFile=-/etc/default/ili9488-daemon
ExecStart=/usr/bin/ili9488-daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
